# deploy/opsledger.service
[Unit]
Description=OpsLedger
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=opsledger
Group=opsledger
EnvironmentFile=/etc/opsledger/opsledger.env
ExecStart=/opt/opsledger/opsledger
Restart=on-failure
RestartSec=5s
WorkingDirectory=/var/lib/opsledger
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/opsledger

[Install]
WantedBy=multi-user.target
